# ops/systemd/snippets-dev.service
[Unit]
Description=Snippets dev server (Vite + API)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/projects/snippets
ExecStart=/usr/bin/env bash -lc 'npm run dev'
Restart=always
RestartSec=2
KillSignal=SIGTERM
TimeoutStopSec=20
Environment=NODE_ENV=development

[Install]
WantedBy=default.target
